# supervisor-crond.ini
; Menjalankan busybox crond di bawah supervisor (bersama php-fpm & nginx).
; crond memproses /etc/crontabs/root yang berisi jadwal generate QR.
[program:crond]
command=crond -f -d 8 -L /dev/stdout
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
