# supervisord.conf
[program:hit_api]
command=/home/hit_api.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/hit_api.log
stderr_logfile=/var/log/hit_api.err.log
